# nginx/files/crontab
#m h dom mon dow user command
*  * 1   *   *   root certbot renew --post-hook "nginx -s reload" >> /var/log/cron.log 2>&1
